# Alpine/testCpu1_service
#!/sbin/openrc-run

command=/usr/bin/php
command_args=/home/test_cgroups/charge_cpu.php